# ops/systemd/greatgo-weekly-report-daemon.service
[Unit]
Description=GreatGo Report Worker Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/greatgo
ExecStart=/usr/bin/npm run report:weekly:daemon
Restart=always
RestartSec=15
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
